# telegramBot.service
[Unit]
Description=Dify Bazi Bot Service (Rust)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi/DifyBaziWorkflow
ExecStart=/home/dietpi/DifyBaziWorkflow/target/release/dify-telegram-bot
Restart=always
RestartSec=5
# Environment file for TOKEN and DIFY_WEBHOOK_URL
EnvironmentFile=/home/dietpi/DifyBaziWorkflow/.env

[Install]
WantedBy=multi-user.target